# etc/profile.d/grep.sh
#!/bin/bash

if [[ -z $GREP_OPTIONS ]]; then                                                      
    export GREP_OPTIONS='--color=always'                                             
fi

