# etc/profile.d/grep.sh
#!/usr/bin/env bash

if [[ -z $GREP_OPTIONS ]]; then
    # Don't use this code.
    #export GREP_OPTIONS='--color=always'

    export GREP_OPTIONS=
fi

